# latest/overlay-image-tools/etc/init.d/scw-gen-machine-id
#!/sbin/openrc-run
# Copyright (c) 2018 Scaleway <opensource@scaleway.com>
# Released under the MIT license.

description="Generate VM UUID"

depend() {
    need root
    after net
}

start() {
    ebegin "Generating machine UUID"
    /usr/local/sbin/scw-gen-machine-id
    eend $?
}
